# etc/supervisor/conf.d/sshd.conf
[program:sshd]
priority=90
command=/usr/sbin/sshd
user=balloon
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true

